# supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
command=/opt/start.sh
